# downloader/config/crontab
*/5 * * * * /usr/bin/xargs -a /data/playlists.txt -n 1 -P 4 /app/download.sh